# deploy/construneves-painel-api.service
[Unit]
Description=Construneves — API Firebird (interna, 127.0.0.1:8091)
After=network.target

[Service]
Type=simple
User=tactical
Group=tactical
WorkingDirectory=/var/www/construneves/painel-api
EnvironmentFile=/var/www/construneves/.env
ExecStart=/var/www/construneves/painel-api/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8091
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
